# docker/supervisor.d/registry.ini
[program:registry]
command=/bin/registry serve /etc/docker/registry/config.yml
priority=500
user=root

